# silver-tier/schedule/windows_task_gmail_watcher.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler XML — Gmail Watcher
  Import with: schtasks /create /xml "windows_task_gmail_watcher.xml" /tn "AIEmployee\GmailWatcher"
  Or: Import via Task Scheduler GUI → Action → Import Task

  IMPORTANT: Update <Command> and <Arguments> paths before importing.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-02-25T00:00:00</Date>
    <Author>AI Employee</Author>
    <Description>AI Employee Gmail Watcher — monitors apple379tree@gmail.com for new emails</Description>
    <URI>\AIEmployee\GmailWatcher</URI>
  </RegistrationInfo>
  <Triggers>
    <!-- Run at system startup -->
    <BootTrigger>
      <Enabled>true</Enabled>
      <Delay>PT2M</Delay>
    </BootTrigger>
    <!-- Also run at logon -->
    <LogonTrigger>
      <Enabled>true</Enabled>
      <Delay>PT1M</Delay>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <RestartOnFailure>
      <Interval>PT5M</Interval>
      <Count>10</Count>
    </RestartOnFailure>
    <Enabled>true</Enabled>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <!-- UPDATE THIS PATH to your Python executable -->
      <Command>C:\Python313\python.exe</Command>
      <!-- UPDATE THIS PATH to match your project location -->
      <Arguments>E:\hackathon-0\silver-tier\watchers\gmail_watcher.py --vault "E:\hackathon-0\silver-tier\AI_Employee_Vault" --interval 120</Arguments>
      <WorkingDirectory>E:\hackathon-0\silver-tier</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
