# silver-tier/schedule/windows_task_orchestrator.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler XML — Orchestrator (runs every 5 minutes)
  Import with: schtasks /create /xml "windows_task_orchestrator.xml" /tn "AIEmployee\Orchestrator"
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-02-25T00:00:00</Date>
    <Author>AI Employee</Author>
    <Description>AI Employee Orchestrator — polls vault every 5 minutes and triggers Claude reasoning loop</Description>
    <URI>\AIEmployee\Orchestrator</URI>
  </RegistrationInfo>
  <Triggers>
    <!-- Repeat every 5 minutes, indefinitely -->
    <TimeTrigger>
      <Repetition>
        <Interval>PT5M</Interval>
        <Duration>P9999D</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2026-02-25T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
    </TimeTrigger>
    <!-- Also at startup -->
    <BootTrigger>
      <Enabled>true</Enabled>
      <Delay>PT3M</Delay>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <ExecutionTimeLimit>PT4M</ExecutionTimeLimit>
    <Enabled>true</Enabled>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Python313\python.exe</Command>
      <Arguments>E:\hackathon-0\silver-tier\watchers\orchestrator.py --once --vault "E:\hackathon-0\silver-tier\AI_Employee_Vault"</Arguments>
      <WorkingDirectory>E:\hackathon-0\silver-tier</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
